# servers/config/blockhttp.conf

start on runlevel [2345]
stop on runlevel [!2345]

setuid vagrant
setgid www-data

env PATH=/home/vagrant/henv/bin
chdir /var/www/http
exec uwsgi --ini app.ini
